# packaging/rpm/99-usbypass.rules
# USBYPASS — route USB mass-storage add/remove events to our handler.
# RPM package variant: shims live in /usr/libexec.

ACTION=="add",    SUBSYSTEM=="block", KERNEL=="sd[a-z][0-9]", ENV{ID_BUS}=="usb", RUN+="/usr/libexec/usbypass-udev-handler add %k"
ACTION=="change", SUBSYSTEM=="block", KERNEL=="sd[a-z][0-9]", ENV{ID_BUS}=="usb", RUN+="/usr/libexec/usbypass-udev-handler add %k"
ACTION=="remove", SUBSYSTEM=="block", KERNEL=="sd[a-z][0-9]",                    RUN+="/usr/libexec/usbypass-udev-handler remove %k"
